# cmake/deb/postinst
deb-systemd-helper disable 'v4l2-stream.service' >/dev/null || true
# if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
# 	deb-systemd-helper unmask 'v4l2-stream.service' >/dev/null || true
# 
# 	# was-enabled defaults to true, so new installations run enable.
# 	if deb-systemd-helper --quiet was-enabled 'v4l2-stream.service'; then
# 		# Enables the unit on first installation, creates new
# 		# symlinks on upgrades if the unit file has changed.
# 		deb-systemd-helper enable 'v4l2-stream.service' >/dev/null || true
# 	else
# 		# Update the statefile to add new symlinks (if any), which need to be
# 		# cleaned up on purge. Also remove old symlinks.
# 		deb-systemd-helper update-state 'v4l2-stream.service' >/dev/null || true
# 	fi
# fi
# if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
# 	if [ -d /run/systemd/system ]; then
# 		systemctl --system daemon-reload >/dev/null || true
# 		if [ -n "$2" ]; then
# 			_dh_action=restart
# 		else
# 			_dh_action=start
# 		fi
# 		deb-systemd-invoke $_dh_action 'v4l2-stream.service' >/dev/null || true
# 	fi
# fi
